# learn-shell/roboshop-shell/cart.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis.roboshopk8.online
Environment=CATALOGUE_HOST=catalogue.roboshopk8.online
Environment=CATALOGUE_PORT=8080
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target